# etc/init.d/pbr-cust
#!/bin/sh /etc/rc.common
# CustRule PBR 自動套用腳本，保證開機及 wg2 ready 後執行

START=99
STOP=10

LOG_FILE="/tmp/pbr-cust.log"

log() {
    UTC_TS=$(date -u +%s)
    TW_TS=$((UTC_TS + 8*3600))
    TW_DATE=$(date -u -d "@$TW_TS" +"%Y-%m-%d %H:%M:%S")
    echo "$TW_DATE - $1" >> $LOG_FILE
}

wait_for_wg2() {
    # 等 wg2 interface up
    for i in $(seq 1 30); do
        ip link show wg2 >/dev/null 2>&1 && return 0
        sleep 1
    done
    return 1
}

cleanup_old() {
    log "Cleaning up CustRule PBR..."
    for tbl in $(seq 1000 4000); do
        ip route flush table $tbl 2>/dev/null
        ip rule | grep "lookup $tbl" | while IFS= read -r line; do
            idx=$(echo "$line" | awk -F: '{print $1}' | tr -d ' ')
            [ -n "$idx" ] && ip rule del pref $idx 2>/dev/null
            log "Deleted ip rule #$idx for table $tbl"
        done
    done
    log "CustRule PBR cleanup complete."
}


cleanup() {
    log "Cleaning up CustRule PBR..."
    
    # 步驟 1: 一次性解析所有需要刪除的規則，並刪除 ip rule
    # 運行一次 ip rule，並讓 awk 處理所有邏輯，避免數千次 grep/awk/tr 呼叫
    ip rule | awk -F: '
        /lookup/ { 
            # 檢查 lookup table 是否在 1000 到 4000 之間
            split($0, arr, "lookup "); 
            table_id = arr[2];
            
            if (table_id >= 1000 && table_id <= 4000) {
                # 提取規則編號 (例如：0: from all)
                idx = $1;
                
                # 執行刪除指令 (使用 system() 執行外部命令)
                # 使用絕對路徑確保 Crontab 環境兼容
                system("/sbin/ip rule del pref " idx " 2>/dev/null");
                
                # 記錄日誌
                # system("/usr/bin/logger -t pbr_cleanup \"Deleted ip rule #" idx " for table " table_id "\"");
                
                # 記錄到 stdout (用於 Crontab 日誌)
                print "Deleted ip rule #" idx " for table " table_id
            }
        }
    '

    # 步驟 2: 清空路由表
    # 雖然仍需 4100 次呼叫，但這是 ip route flush 的最低要求
    for tbl in $(/usr/bin/seq 1000 4000); do
        /sbin/ip route flush table $tbl 2>/dev/null
    done
    
    log "CustRule PBR cleanup complete."
}


apply_rules() {
    log "Applying CustRule PBR (diff mode)..."

    # 從 uci 算出「應有」的規則: "src_ip table_id dev" 列表
    WANT=$(uci show pbr | grep "=policy" | while IFS= read -r line; do
        idx=$(echo "$line" | cut -d'=' -f1)
        name=$(uci get "$idx.name" 2>/dev/null)
        case "$name" in CustRule*) ;; *) continue ;; esac
        enabled=$(uci get "$idx.enabled" 2>/dev/null)
        [ "$enabled" = "0" ] && continue
        src=$(uci get "$idx.src_addr" 2>/dev/null)
        dev=$(uci get "$idx.dest_addr" 2>/dev/null)
        [ -z "$src" ] || [ -z "$dev" ] && continue
        table_id=$(echo "$name" | tr -dc '0-9')
        echo "$src $table_id $dev"
    done)

    # 從現有 ip rule 讀出「現有」的規則: "src_ip table_id"
    HAVE=$(ip rule list | awk '$1=="200:" && $2=="from" && $3!="all" {
        tbl=$5+0; if(tbl>=1000 && tbl<=4000) print $3, tbl}')

    # 要刪的 = HAVE 有但 WANT 沒有
    echo "$HAVE" | while read _src _tbl; do
        [ -z "$_src" ] && continue
        if ! echo "$WANT" | awk '{print $1, $2}' | grep -q "^${_src} ${_tbl}$"; then
            ip rule del prio 200 from "$_src" lookup "$_tbl" 2>/dev/null
            ip route flush table "$_tbl" 2>/dev/null
            log "Removed rule: $_src lookup $_tbl"
        fi
    done

    # 要加的 = WANT 有但 HAVE 沒有（或 route 需更新）
    echo "$WANT" | while read _src _tbl _dev; do
        [ -z "$_src" ] && continue
        _have_rule=$(echo "$HAVE" | awk -v s="$_src" -v t="$_tbl" '$1==s && $2==t')

        if [ -z "$_have_rule" ]; then
            ip rule add prio 200 from "$_src" lookup "$_tbl" 2>/dev/null
            log "Added rule: $_src lookup $_tbl"
        fi

        # 確保 route 正確（dev 可能有變）
        _cur_route=$(ip route show table "$_tbl" 2>/dev/null | grep '^default')
        case "$_dev" in
            wan)
                WAN_GW=$(ip route show default dev wan 2>/dev/null | awk '{print $3; exit}')
                _want_route="default via $WAN_GW dev wan"
                ;;
            wg*)
                _want_route="default dev $_dev"
                ;;
            *)
                log "Warning: unknown device $_dev for $_src, skip"
                continue
                ;;
        esac

        if [ "$_cur_route" != "$_want_route" ]; then
            ip route flush table "$_tbl" 2>/dev/null
            case "$_dev" in
                wan) ip route add default via "$WAN_GW" dev wan table "$_tbl" 2>/dev/null ;;
                wg*) ip route add default dev "$_dev" table "$_tbl" 2>/dev/null ;;
            esac
            log "Updated route: $_src → $_dev (table $_tbl)"
        fi
    done

    ip route flush cache
    log "CustRule PBR applied successfully."

    # 更新 check-pbr-wg cache（供 DOWN/UP 無感切換用）
    # fwmark cache: 每個 wg 介面的 fwmark ip rule
    ip rule list | awk '/fwmark.*lookup pbr_wg/{
        prio=$1+0
        for(i=1;i<=NF;i++) if($i=="fwmark") fm=$(i+1)
        tbl=$NF
        gsub(/^pbr_/,"",tbl)
        print prio, fm, "pbr_" tbl > "/tmp/check-pbr-wg." tbl ".rule"
    }'

    # custrule cache: 每個 wg 介面對應的 per-IP rule 列表
    # 先清掉舊 cache
    rm -f /tmp/check-pbr-wg.*.custrules 2>/dev/null
    # 從 uci 重建（包含 enabled=0 的也記錄，DOWN 時才知道要刪哪條）
    uci show pbr | grep "=policy" | while IFS= read -r line; do
        _idx=$(echo "$line" | cut -d'=' -f1)
        _name=$(uci get "$_idx.name" 2>/dev/null)
        case "$_name" in CustRule*) ;; *) continue ;; esac
        _src=$(uci get "$_idx.src_addr" 2>/dev/null)
        _dev=$(uci get "$_idx.dest_addr" 2>/dev/null)
        [ -z "$_src" ] || [ -z "$_dev" ] && continue
        _tbl=$(echo "$_name" | tr -dc '0-9')
        echo "$_src $_tbl" >> "/tmp/check-pbr-wg.${_dev}.custrules"
    done
    log "Cache updated (fwmark + custrules)"
}

start() {
    log "Starting CustRule PBR auto..."

    # ..................
    #cp /etc/init.d/pbr-cust /etc/myscript/pbr-cust
    #cp /etc/hotplug.d/iface/99-pbr-cust /etc/myscript/99-pbr-cust
    #log "Copied script to /etc/config/pbr-cust  /etc/myscript/99-pbr-cust"

    

    if wait_for_wg2; then
        apply_rules
    else
        log "Warning: wg2 interface not found, skipped apply"
    fi
}

stop() {
    log "Stopping CustRule PBR auto..."
    cleanup
}







# chmod +x /etc/init.d/pbr-cust
# /etc/init.d/pbr-cust start
# /etc/init.d/pbr-cust enable

# chmod +x /etc/hotplug.d/iface/99-pbr-cust

